# docker/supervisor.conf.d/janitoo_hostsensor.conf
[program:jnt_hostsensor]
user=root
autostart=true
autorestart=true
command=jnt_hostsensor -c /opt/janitoo/etc/janitoo_hostsensor.conf front
